# bottlenecks/src/api/supervisor/conf.d/system_monitor.conf
[program:system_monitor]
command=python3 /app/monitor/system_monitor.py
autostart=true
autorestart=true
user=vagrant
redirect_stderr=true
stdout_logfile=/var/log/%(program_name)s